# crontab
44 4 * * 6 /usr/bin/docker exec -t nginx nginx -s reload
